# docker/rootfs/etc/services.d/nginx/finish
XSym
0009
83f4a1cefd1c7cb21e72045eb25c914b
/bin/true
                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                                      